# etc/hotplug.d/iface/00-connect
[ "ifup" = "${ACTION}" ] && [ "wlan0" = "${DEVICE}" ] && [ "wwan" = "${INTERFACE}" ] && {

	python /root/set_connected.py
}
